# contrib/playground/deb/DEBIAN/postinst
systemctl daemon-reload
systemctl enable ufw --now

ufw allow 12321/tcp
ufw reload

systemctl enable eyot-playground.service --now