# contrib/playground/deb/DEBIAN/postinst
systemctl daemon-reload
systemctl enable ufw --now

systemctl enable eyot-playground.service --now