# contrib/playground/deb/DEBIAN/postinst
systemctl daemon-reload
systemctl enable eyot-playground.service --now